# scripts/systemd/bbs-data-to-gcp.timer
[Unit]
Description=timer bbs.jpnkn.com data backup to GCP

[Timer]
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
